# infrastructure/peladn-host/backup-scripts/backup-metrics-textfile.service
# /etc/systemd/system/backup-metrics-textfile.service
[Unit]
Description=Emit homelab backup freshness/size metrics for node-exporter
[Service]
Type=oneshot
ExecStart=/opt/backup-scripts/backup-metrics-textfile.sh
